# smstools3-openwrt/files/smstools3
#!/bin/sh /etc/rc.common
# "new(er)" style init script
# Look at /lib/functions/service.sh on a running system for explanations of what other SERVICE_
# options you can use, and when you might want them.

START=80
APP=smstools3
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1

start() {
        mkdir -p /var/smstools3/sms/outgoing
		mkdir -p /var/smstools3/sms/incoming
		mkdir -p /var/smstools3/sms/sent
		mkdir -p /var/smstools3/sms/checked
		service_start /usr/sbin/smsd
}

stop() {
        service_stop /usr/sbin/smsd
}
