# dragino-yun-scripts/files/etc/init.d/rngd-turn-off
#!/bin/sh /etc/rc.common

START=92

start() {
	/etc/init.d/rngd stop
}

