# globacom-pkg/files/etc/uci-defaults/96_globacom_default
#!/bin/sh

#rename-hostname
uci set firewall.@forwarding[0].dest=lan
uci commit firewall
/etc/init.d/firewall restart
exit 0