# auto-update-packages/files/etc/init.d/update_pkg
#!/bin/sh /etc/rc.common

#A script to update selected packages. 
## Copyright (c) 2015 Dragino Tech <support@dragino.com>

START=99

start()
{
	# Configure auto update periodically.
	# Remove previous configure
	sed '/update_packages/d' /etc/crontabs/root -i
	UGD_ON_TIME=`uci get auto_update_pkg.general.update_on_time`
	if [ $UGD_ON_TIME == '1' ];then
		MIN=`uci get auto_update_pkg.general.min`
		HOUR=`uci get auto_update_pkg.general.hour`
		DATE_OF_MONTH=`uci get auto_update_pkg.general.date_of_month`
		MONTH=`uci get auto_update_pkg.general.month`
		DATE_OF_WEEK=`uci get auto_update_pkg.general.date_of_week`
		echo "$MIN $HOUR $DATE_OF_MONTH $MONTH $DATE_OF_WEEK /usr/bin/update_packages.sh" >>  /etc/crontabs/root
	fi 
	

	# Check if need to check update on boot.
	UGD_ON_BOOT=`uci get auto_update_pkg.general.update_on_boot`
	if [ $UGD_ON_BOOT == '1' ];then
		/usr/bin/update_packages.sh &
	fi
}