# dragino2/files/etc/uci-defaults/97_yun_default
#!/bin/sh
rm /etc/avahi/services/http.service

exit 0